# rootfs/etc/periodic/daily/acmetool-reconcile
#!/usr/bin/with-contenv sh
/usr/local/bin/acmetool --batch reconcile
